# kali-autonomous/systemd/kali-autonomous.timer
[Unit]
Description=Run Kali Autonomous Security Workflow repeatedly

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Unit=kali-autonomous.service
Persistent=true

[Install]
WantedBy=timers.target
